# archive/CQT Repo/repos/logger/doread.cron
# file that reads out the ohmmeter every minute
* * * * * /home/qitlab/programs/logger/doread
